# debian/vhostmd.postinst
#!/bin/sh

systemctl enable --now vhostmd

#DEBHELPER#
exit 0
